# qcontrold.socket
[Socket]
ListenStream=/run/qcontrol.sock

[Install]
WantedBy=sockets.target
